# certwarden-ha-client/rootfs/etc/s6-overlay/s6-rc.d/certwarden-ha-client/run
#!/usr/bin/with-contenv bashio

# Start the service
exec /usr/bin/certwarden-ha-client